# rpigw.service
[Unit]
Description=Raspberry Pi IQRF and GSM gateway

[Service]
ExecStart=rpigw
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=rpigw.service
